# root/etc/services.d/sinatra/run
#!/usr/bin/with-contenv bash

cd /opt/dockerhub-stats/
ruby stats.rb
